# simplified/tasks/scheduled_tasks/file maintenance/_task_defs/_start_new_batch_run_session.xml
<?xml version="1.0" encoding="UTF-16"?>
    <Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
      <RegistrationInfo>
        <Version>             1                                                                        </Version>
        <Source>              N/A                                                      </Source>
        <Date>                2024-04-04T14:03:58.8742243                                                         </Date>
        <Author>              DSKTP-HOME-JEFF\jeffers                                                  </Author>
        <Description>         Part of FilmCab, set up the batch run session for tracking flow. Nothing should run if there isn't an active record in the batch_run_session table. . # 1</Description>
        <URI>                 \FilmCab\file maintenance\_start_new_batch_run_session                                                                     </URI>
      </RegistrationInfo>
      <Triggers>
                  <CalendarTrigger id="start_batch_run_session_on_time">
          
          <StartBoundary>         2024-04-04T00:01:00.0000000    </StartBoundary>
          <ExecutionTimeLimit>    PT72H                 </ExecutionTimeLimit>
          <Enabled>               true                  </Enabled>
          <ScheduleByDay>
            <DaysInterval> 1 </DaysInterval>
          </ScheduleByDay>
        </CalendarTrigger>
      </Triggers>
      <Principals>
        <Principal id="Author">
          <UserId>S-1-5-21-260979430-3554011381-420227292-1001</UserId>
          <LogonType>Password</LogonType>
          <RunLevel>HighestAvailable</RunLevel>
        </Principal>
      </Principals>
      <Settings>
        <MultipleInstancesPolicy>    IgnoreNew     </MultipleInstancesPolicy>
        <DisallowStartIfOnBatteries> true          </DisallowStartIfOnBatteries>
        <StopIfGoingOnBatteries>     true          </StopIfGoingOnBatteries>
        <AllowHardTerminate>         true          </AllowHardTerminate>
        <StartWhenAvailable>         false         </StartWhenAvailable>
        <RunOnlyIfNetworkAvailable>  false         </RunOnlyIfNetworkAvailable>
        <IdleSettings>
          <StopOnIdleEnd>            true          </StopOnIdleEnd>
          <RestartOnIdle>            false         </RestartOnIdle>
          <WaitTimeout>              PT1H          </WaitTimeout>
          <Duration>                 PT1M          </Duration>
        </IdleSettings>
        <AllowStartOnDemand>               true                      </AllowStartOnDemand>
        <Enabled>                          true                      </Enabled>
        <Hidden>                           false                     </Hidden>
        <RunOnlyIfIdle>                    false                     </RunOnlyIfIdle>
        <DisallowStartOnRemoteAppSession>  false                     </DisallowStartOnRemoteAppSession>
        <UseUnifiedSchedulingEngine>       true                      </UseUnifiedSchedulingEngine>
        <WakeToRun>                        true                      </WakeToRun>
        <ExecutionTimeLimit>               PT2H     </ExecutionTimeLimit>
        <Priority>                         7                         </Priority>
      </Settings>
      <Actions Context="Author">
        <Exec>
          <Command>"C:\Program Files\PowerShell\7-preview\pwsh.exe"</Command>
          <Arguments>-ExecutionPolicy Bypass -Command ". 'D:\qt_projects\FilmCab\simplified\tasks\scheduled_tasks\file maintenance\_start_new_batch_run_session.ps1'; exit $LASTEXITCODE"</Arguments>
          <WorkingDirectory>D:\qt_projects\filmcab</WorkingDirectory>
        </Exec>
      </Actions>
    </Task>    
